# launchd/com.github.Coco422.doubao-voice-wetype-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROGRAM__</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_LAUNCHD_LABEL</key>
    <string>__LABEL__</string>
    <key>RESTORE_IME_ID</key>
    <string>com.tencent.inputmethod.wetype.pinyin</string>
    <key>VOICE_IME_ID</key>
    <string>com.bytedance.inputmethod.doubaoime.pinyin</string>
    <key>VOICE_IME_ALIASES</key>
    <string>com.bytedance.inputmethod.doubaoime</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/doubao-voice-wetype-agent.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/doubao-voice-wetype-agent.stderr.log</string>

  <key>WorkingDirectory</key>
  <string>__HOME__</string>
</dict>
</plist>
